# supervisord.conf
[supervisord]
nodaemon=true

[program:couchdb]
command=/usr/local/bin/couchdb -p /usr/local/var/run/couchdb/couchdb.pid
autorestart=true
autostart=true