# opt/supervisord.conf
[supervisord]
nodaemon=true

[program:couchdb]
command=couchdb -b -p /usr/local/var/run/couchdb/couchdb.pid
user=couchdb
directory=/usr/local/var/lib/couchdb
autorestart=true
autostart=true