# opt/supervisord.conf
[supervisord]
nodaemon=true

[program:couchdb]
command=couchdb
user=couchdb
directory=/usr/local/var/lib/couchdb
autorestart=true
autostart=true